# root/etc/s6-overlay/s6-rc.d/svc-mod-beets-sync-mod/run
#!/usr/bin/with-contenv bash

exec gosu "$PUID:$PGID" /etc/s6-overlay/s6-rc.d/svc-mod-beets-sync-mod/main "$@"
